# deploy/systemd/expense-voice-fastapi.service
[Unit]
Description=Expense Voice FastAPI Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/familybudget/python
Environment=PYTHONUNBUFFERED=1
Environment=EXPENSE_VOICE_MODEL=small
Environment=EXPENSE_VOICE_FASTAPI_TOKEN=
ExecStart=/var/www/familybudget/python/.venv/bin/uvicorn expense_voice.fastapi_app:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
